# crontab
*/5 * * * * c:/php5/php d:/patch/ping.php